# resources/crontab
#
5 14 * * 5 open "x-apple.systempreferences:com.apple.BluetoothSettings"
#
